# crontab
* * * * * /dstar/scripts/dstarsms.pl >/dstar/scripts/dstarsms.cron.log 2>&1
